# datebase/postgresql/supervisor/conf.d/zpcustom.conf
[supervisord]
nodaemon=true
logfile_maxbytes=50MB 
logfile_backups=10

[program:postgresql]
command = postgres
startsecs=3
stopwaitsecs=3
autostart=true
autorestart=true


