# datebase/redis/supervisor/conf.d/zpcustom.conf
[supervisord]
nodaemon=true
logfile_maxbytes=50MB 
logfile_backups=10

[program:redis]
command = redis-server /etc/redis/redis.conf
startsecs=3
stopwaitsecs=3
autostart=true
autorestart=true


